# 04-configurando_ec2_com_user_data/files/user_data
#cloud-config
write_files:
  - path: /etc/init.d/custom_startup
    permissions: "0755"
    owner: root
    content: |
      #!/bin/sh -
      # Custom startup of applications
      #
      ### BEGIN INIT INFO
      # Provides: Custom startup of applications
      # Default-Start: 2 3 4 5
      # Default-Stop: 0 1 6
      # Required-Start: $network $remote_fs
      # Required-Stop: $remote_fs $network
      # Short-Description: That's  all folks
      # Description: That's  all folks
      ### END INIT INFO

      for i in {1..20}
      do
        docker info && break
        sleep 1
      done

      sudo docker run -d -p 3000:3000 grafana/grafana